# systemd/update_do_blocklist.service
[Unit]
Description=Update DigitalOcean blocklist and import into CrowdSec (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# This service runs as root, which is required to interact with the cscli.
User=root
# Adjust ExecStart arguments as desired.
# The ExecStart path below is a placeholder and will be replaced by the systemd_install script.
ExecStart=/usr/local/bin/update_DO_blocklist -k --no-sudo -f /var/cache/do_blocklist.csv sync
# Write logs to the system journal; you can also capture to a file with a wrapper or use StandardOutput=append:/var/log/do_blocklist.log
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
